# scripts/com.ask149.aidaemon.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ask149.aidaemon.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ashishkshirsagar/Projects/active/aidaemon/scripts/watchdog.sh</string>
    </array>

    <!-- Run every 1800 seconds (30 minutes) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Also run immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Environment for node/npx (MCP servers need this) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ashishkshirsagar</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ashishkshirsagar/.config/aidaemon/data/logs/watchdog-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ashishkshirsagar/.config/aidaemon/data/logs/watchdog-launchd.log</string>
</dict>
</plist>
